# backups/etc/systemd/system/backups.service
[Unit]
Description=Makes backup of /data
RequiresMountsFor=/data
ConditionACPower=true

[Service]
Type=oneshot
ExecStart=/root/bin/backups.sh
